# deploy/thereview-api.service
[Unit]
Description=The Review API
After=network.target postgresql.service

[Service]
Type=simple
User=thereview
Group=thereview
WorkingDirectory=/opt/thereview/backend
Environment=PATH=/opt/thereview/backend/.venv/bin:/usr/bin
ExecStart=/opt/thereview/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
